# fff-ap@.service
[Unit]
After=network.target

[Service]
ExecStart=/home/pi/fff-ap/fff-ap.sh %I

[Install]
WantedBy=multi-user.target
